# roles/backend/templates/backend.service
[Unit]
Description=Backend Service

[Service]
User=sri
Environment=DB_HOST="{{ lookup('aws_ssm', 'expense_backend_mysql_ip', region='us-east-1') }}"
ExecStart=/bin/node /app/index.js
SyslogIdentifier=Backend

[Install]
WantedBy=multi-user.target